# ap0-watchdog.timer
[Unit]
Description=Run ap0 watchdog every 30 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=30

[Install]
WantedBy=timers.target