# ap0-watchdog.service
[Unit]
Description=Check and restore ap0 if missing
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-ap0.sh